# board/common/rootfs/etc/profile.d/convenience.sh
alias vim='vi'
alias emacs='mg'
alias edit='mg'
alias editor=edit
alias sensible-editor=edit
alias sensible-pager=less
