# board/common/rootfs/etc/rc.local
#!/bin/sh

# alt: sun12x22 -d
#setfont cybercafe -d
#loadkeys sv-latin1

