# initramfs/etc/inittab
::sysinit:/etc/init.d/rcS
::askfirst:-/bin/sh

# Stuff to do when restarting the init process
::restart:/sbin/init

# Stuff to do before rebooting
::ctrlaltdel:/sbin/reboot
::shutdown:/etc/shutdown

::once:/sbin/watchdog -T 16 -t 2 -F

::respawn:/sbin/inetd -f
#::respawn:sh -c "sleep 3; exec /sbin/udhcpc -i eth0 -p /var/run/udhcpc.pid -f"
::respawn:/sbin/ifplugd -n -a -i eth0 -I -l

::respawn:/sbin/httpd -r "Administrator Login" -f

::respawn:/sbin/getty 115200 ttyGS0 vt100
