# legacy/system/skeleton/etc/rc.local
#!/bin/sh

echo "rc.local has run" >/tmp/hello-rc.local
